# profiles/templates/3.1/3_ac_desktop_desktop/2-desktop-user/60-checkshares-0.1/.config/autostart/checkshares.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=checkshares
Comment=Bash script for check domain shares availability, mount via automount and add gtk bookmarks.
Exec=/usr/bin/checkshares.sh
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false

